# bbs.conf
[program:bbs]
command=/usr/local/bin/gunicorn wsgi -c gunicorn.config.py
directory=/root/web_hannah
autostart=true
autorestart=true
