# mongodb.conf
[program:mongodb]
command=/usr/bin/mongod --dbpath /root/web_hannah/data/db --port 27023 --logpath /root/web_hannah/log/mongod.log 
directory=/root/web_hannah
autostart=true
autorestart=true
